# apsystems-ez1-homa/rootfs/etc/services.d/apsystems-ez1/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "apsystems-ez1 service starting ..."

# Fetch MQTT settings
source get_mqtt_settings.sh

# Activate the virtual environment
source .venv/bin/activate

# Fetch Home Assistant config to get latitude and longitude
# needs the "homeassistant_api: true" and "hassio_api: true" permissions in config.yaml
api_config=$(curl -s -H "Authorization: Bearer $SUPERVISOR_TOKEN" http://supervisor/core/api/config)

# create the service configuration
export APS_ECU_IP="$(bashio::config 'aps_ecu_ip')"
export APS_ECU_PORT=$(bashio::config 'aps_ecu_port' '8050')
export APS_ECU_UPDATE_INTERVAL=$(bashio::config 'aps_ecu_update_interval' '15')
export APS_ECU_STOP_AT_NIGHT=$(bashio::config 'aps_ecu_stop_at_night' 'false')
# export APS_ECU_POSITION_LAT=$(bashio::info.latitude) # bashio::info.latitude is not supported any more
# export APS_ECU_POSITION_LNG=$(bashio::info.longitude) # use the API to get the latitude and longitude instead
export APS_ECU_POSITION_LAT=$(echo "$api_config" | jq -r '.latitude')
export APS_ECU_POSITION_LNG=$(echo "$api_config" | jq -r '.longitude')
export APS_ECU_TIMEZONE="$(bashio::info.timezone)"

export MQTT_BROKER_HOST="${MQTT_HOST}"
export MQTT_BROKER_PORT="${MQTT_PORT}"
export MQTT_BROKER_USER="${MQTT_USER}"
export MQTT_BROKER_PASSWD="${MQTT_PASSWORD}"
export MQTT_BROKER_SECURED_CONNECTION=False
export MQTT_CLIENT_ID="APsystemsEZ1mqtt_addon"
export HOMA_ENABLED=True
export HOMA_ROOM="$(bashio::config 'homa_room')"
export HOMA_NAME="$(bashio::config 'device_name')"
export HASS_ENABLED=True
export HASS_DEVICE_NAME="$(bashio::config 'device_name')"
export HASS_NAME_PREFIX="Solar "
export HASS_AREA="$(bashio::config 'hass_area')"

if bashio::config.true 'debug'; then
    bashio::log.info "APS_ECU_POSITION_LAT: $APS_ECU_POSITION_LAT"
    bashio::log.info "APS_ECU_POSITION_LNG: $APS_ECU_POSITION_LNG"
    bashio::log.info "APS_ECU_TIMEZONE: $APS_ECU_TIMEZONE"
fi

bashio::config.true 'debug' && extra_parms="-d" || extra_parms=""

# Start the apsystems-ez1 service
exec python3 -m apsystems_ez1_mqtt $extra_parms

#tail -f /dev/null
